# frpcDeploy/frpc/systemd/frpc@.service
[Unit]
Description=Frp Client Service
After=network.target

[Service]
Type=simple
User=<user>
Restart=on-failure
RestartSec=5s
ExecStart=<dir>/frpc -c <dir>/%i.ini
ExecReload=<dir>/frpc reload -c <dir>/%i.ini
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
